# docker/xenial/bootable/smartd.service.d/0no-virt.conf
[Unit]
|ConditionKernelCommandLine=placeholder
|ConditionVirtualization=false
