# docker/xenial/xenhv/systemd-system-dnsmasq.service.d/local.conf
[Service]
RestartSec=1s
Restart=on-failure
