# docker/xenial/installer/all-console-setup/addcons.service
[Unit]
Description=serial-console for all 16550A-serial ports

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/startcons.sh
