# docker/xenial/bootable/platform-info/platform-info.service
[Unit]
Description=platform information service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/platform-info.sh
